# services/my-timers/sysinfo.service
[Unit]\nDescription=Generate sysinfo JSON for dashboard\n\n[Service]\nType=oneshot\nExecStart=/usr/local/bin/sysinfo-json > /var/www/html/vps-dashboard/data/sysinfo.json\n